# taxhub-service.conf
[program:taxhub]
command = APP_PATH/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = APP_PATH/var/log/taxhub-errors.log
redirect_stderr = true
